# scripts/com.cssrs.system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cssrs.system</string>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Developer/cssrs-system/backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/Developer/cssrs-system/.venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
        <string>--log-level</string>
        <string>info</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PROJECT_DIR</key>
        <string>/Users/YOUR_USERNAME/Developer/cssrs-system</string>
        <key>CLOUD_GATEWAY_URL</key>
        <string>http://YOUR_SERVER:8888</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Developer/cssrs-system/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Developer/cssrs-system/logs/launchd.err.log</string>
</dict>
</plist>
