# scripts/cssrs-local.service
[Unit]
Description=C-SSRS Local Service
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/cssrs-system/backend
ExecStart=/path/to/cssrs-system/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitInterval=60

[Install]
WantedBy=multi-user.target
